# mincha-scraper/mincha-scraper.service
[Unit]
Description=Mincha Time Scraper Service
After=network.target

[Service]
Type=oneshot
User=xander
WorkingDirectory=/var/www/trmnl-zmanim/mincha-scraper
ExecStart=/var/www/trmnl-zmanim/mincha-scraper/run_mincha_scraper.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
